# oracle/.config/autostart/register.desktop
[Desktop Entry]
Name=name
Exec=/home/oracle/.config/autostart/register
Type=Application
